# supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
pidfile=/var/run/supervisord/supervisord.pid
childlogdir=/var/log/supervisord/
logfile_maxbytes=5MB                           ; maximum size of logfile before rotation
logfile_backups=3                              ; number of backed-up logfiles
loglevel=error

[program:nextcloud-worker]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=sudo -E -u www-data php occ background-job:worker -t 60 'OC\TaskProcessing\SynchronousBackgroundJob'
autostart=true
autorestart=true
startretries=5
process_name=%(program_name)s_%(process_num)02d
numprocs=%(ENV_WORKER_COUNT)s

[program:cron]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/cron.sh
autostart=true
autorestart=true
startretries=3
